# files/graylog-cookbooks/runit/files/default/graylog-runsvdir.service
[Unit]
Description=Graylog Runit supervision process

[Service]
ExecStart=/opt/graylog/embedded/bin/runsvdir-start
Restart=always
