# files/graylog2-cookbooks/runit/files/default/graylog2-runsvdir.service
[Unit]
Description=Graylog2 Runit supervision process

[Service]
ExecStart=/opt/graylog2/embedded/bin/runsvdir-start
Restart=always
